# deploy/gassigeher.service
[Unit]
Description=Gassigeher Dog Walking Booking System
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=gassigeher
Group=gassigeher
WorkingDirectory=/var/gassigeher
EnvironmentFile=/var/gassigeher/config/.env
ExecStart=/var/gassigeher/bin/gassigeher
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/gassigeher/data /var/gassigeher/uploads /var/gassigeher/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gassigeher

[Install]
WantedBy=multi-user.target
